# server/scholar.service
[Unit]
Description=OpenSource Scholar - Open Access Academia Platform
Documentation=https://github.com/yourusername/opensource-scholar
After=network.target

[Service]
Type=simple
User=scholar
Group=scholar
WorkingDirectory=/opt/scholar/server

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=-/opt/scholar/server/.env

# Start command
ExecStart=/usr/bin/node src/index.js

# Restart policy
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=scholar

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/scholar/server/data

[Install]
WantedBy=multi-user.target
